# feeds/apps/connect_listen/files/connect_listen.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=89

start() {
	/usr/local/bin/connect_listen /dev/ttyS1 11565 11578
}

